# templates/etc/supervisor/conf.d/archive.conf
[program:archive]
command={{ root_prefix }}/var/cnx/venvs/archive/bin/pserve {{ root_prefix }}/etc/cnx/archive/app.ini
