# templates/etc/supervisor/conf.d/authoring.conf
[program:authoring]
command={{ root_prefix }}/var/cnx/venvs/authoring/bin/pserve {{ root_prefix }}/etc/cnx/authoring/app.ini
