# templates/etc/supervisor/conf.d/publishing.conf
[program:publishing]
command={{ root_prefix }}/var/cnx/venvs/publishing/bin/pserve {{ root_prefix }}/etc/cnx/publishing/app.ini
