# templates/etc/supervisor/conf.d/zclient.conf
[program:zclient]
command={{ root_prefix }}/var/lib/cnx/cnx-buildout/bin/instance console
user=www-data
directory={{ root_prefix }}/var/lib/cnx/cnx-buildout
stopwaitsecs=60
