# roles/pdf_gen/templates/etc/supervisor/conf.d/pdf_gen.conf
[program:pdf_gen]
command=/var/lib/cnx/cnx-buildout/bin/pdf_gen console
user=www-data
directory=/var/lib/cnx/cnx-buildout
stopwaitsecs=120
environment=HOME="/var/www",USER="www-data"
