# roles/iptables/files/etc/network/if-up.d/iptables
#!/bin/bash
[[ -e /etc/iptables.up.rules ]] && iptables-restore /etc/iptables.up.rules
