# files/etc/supervisor/conf.d/archive.conf
[program:archive]
command=/usr/local/var/cnx/venvs/archive/bin/pserve /usr/local/etc/cnx/archive/app.ini
