# files/etc/supervisor/conf.d/authoring.conf
[program:authoring]
command=/usr/local/var/cnx/venvs/authoring/bin/pserve /usr/local/etc/cnx/authoring/app.ini
