# files/etc/supervisor/conf.d/publishing.conf
[program:publishing]
command=/usr/local/var/cnx/venvs/publishing/bin/pserve /usr/local/etc/cnx/publishing/app.ini
