# files/etc/supervisor/conf.d/smtp_server.conf
[program:smtp_server]
command=/usr/local/bin/smtp_server.py
