# roles/publishing/templates/etc/supervisor/conf.d/post_publication.conf
[program:post_publication]
command={{ root_prefix }}/var/cnx/venvs/publishing/bin/cnx-publishing-post-publication {{ root_prefix }}/etc/cnx/publishing/app.ini
