# roles/publishing/templates/etc/supervisor/conf.d/post_publication.conf
{% for i in range(0, hostvars[inventory_hostname].post_publication_count|default(1), 1) %}
[program:post_publication{{ i }}]
command=/var/cnx/venvs/publishing/bin/cnx-publishing-post-publication /etc/cnx/publishing/app.ini

{% endfor %}

[group:post_publication]
programs={% for i in range(0, hostvars[inventory_hostname].post_publication_count|default(1), 1) %}post_publication{{ i }}{% if not loop.last %},{% endif %}{% endfor %}
