# roles/publishing/templates/etc/supervisor/conf.d/publishing.conf
{% for i in range(0, hostvars[inventory_hostname].publishing_count|default(1), 1) %}
[program:publishing{{ i }}]
command=/var/cnx/venvs/publishing/bin/pserve --server-name=instance{{ i }} /etc/cnx/publishing/app.ini

{% endfor %}

[group:publishing]
programs={% for i in range(0, hostvars[inventory_hostname].publishing_count|default(1), 1) %}publishing{{ i }}{% if not loop.last %},{% endif %}{% endfor %}
