# roles/zclient/templates/etc/supervisor/conf.d/zclients.conf
{% for i in range(0, hostvars[inventory_hostname].zclient_count|default(1), 1) %}
[program:zclient_instance{{ i }}]
command=/var/lib/cnx/cnx-buildout/bin/instance{{ i }} console
user=www-data
directory=/var/lib/cnx/cnx-buildout
stopwaitsecs=60

{% endfor %}


[program:unoconv]
command=unoconv -l
user=www-data
directory=/var/lib/cnx/cnx-buildout/scripts
stopwaitsecs=60
