# roles/legacy_varnish/templates/etc/cron.daily/copy_varnish_logs
#!/bin/sh
set -e

INSTANCE_NAME="{{ inventory_hostname }}-logs"

rsync -tr --chmod=0700 /var/log/varnish/ /var/www/varnish_logs/$INSTANCE_NAME/
