# roles/archive/files/etc/cron.daily/archive-restart
#!/bin/bash
#exec 1>> /tmp/archive-restart.$$.stdout
#exec 2>> /tmp/archive-restart.$$.stderr

#set -vx

production_dir=/opt/production

cd "$production_dir/bin"

source activate # adds appropriate dir to path and sets virtual env

idle_before=$(psql -U rhaptos repository -tA -c "SELECT COUNT(*) FROM pg_stat_activity WHERE (state = 'idle');")

for instance in "$production_dir"/ca-*
do
    (
        cd "$instance"
        pserve production.ini stop
        sleep 5
        pserve production.ini start
        sleep 5
     )
done

(( idle_after = idle_before ))
(( tries = 0 ))

while (( idle_after >= idle_before && ++tries <= 10 ))
do
    sleep 60
    idle_after=$(psql -U rhaptos repository -tA -c "SELECT COUNT(*) FROM pg_stat_activity WHERE (state = 'idle');")
#    echo "$(date) Try $tries, Before $idle_before, After $idle_after" >&2
done

logger -p local0.notice -t "ARCHIVE-RESTART[$$]" "Postgresql idle connections before restart: $idle_before, after: $idle_after"
