# roles/zeo/templates/etc/supervisor/conf.d/zeo.conf
[program:zeo]
command=/var/lib/cnx/cnx-buildout/bin/zeoserver fg
user=www-data
directory=/var/lib/cnx/cnx-buildout
stopwaitsecs=60
