# roles/archive/templates/etc/supervisor/conf.d/archive.conf
{% for i in range(0, hostvars[inventory_hostname].archive_count|default(1), 1) %}
[program:archive{{ i }}]
command=/var/cnx/venvs/archive/bin/pserve --server-name=instance{{ i }} /etc/cnx/archive/app.ini

{% endfor %}

[group:archive]
programs={% for i in range(0, hostvars[inventory_hostname].archive_count|default(1), 1) %}archive{{ i }}{% if not loop.last %},{% endif %}{% endfor %}
