# roles/authoring/templates/etc/supervisor/conf.d/authoring.conf
{% for i in range(0, hostvars[inventory_hostname].authoring_count|default(1), 1) %}
[program:authoring{{ i }}]
command=/var/cnx/venvs/authoring/bin/pserve --server-name=instance{{ i }} /etc/cnx/authoring/app.ini

{% endfor %}

[group:authoring]
programs={% for i in range(0, hostvars[inventory_hostname].authoring_count|default(1), 1) %}authoring{{ i }}{% if not loop.last %},{% endif %}{% endfor %}
