# roles/publishing_worker/templates/etc/supervisor/conf.d/publishing_celery_workers.conf
{% for i in range(0, hostvars[inventory_hostname].publishing_worker_count|default(1), 1) %}
[program:publishing_celery_worker{{ i }}]
environment=PYRAMID_INI="/etc/cnx/publishing/app.ini"
command=/var/cnx/venvs/publishing/bin/celery worker -A cnxpublishing -Q default -Q deferred
user=www-data

{% endfor %}

[group:publishing]
programs={% for i in range(0, hostvars[inventory_hostname].publishing_worker_count|default(1), 1) %}publishing_celery_worker{{ i }}{% if not loop.last %},{% endif %}{% endfor %}
